# system_files/etc/profile.d/bazzite-neofetch.sh
alias neofetch='/usr/bin/fastfetch -c /usr/share/ublue-os/bazzite/fastfetch.jsonc --logo "$(ls /usr/share/ublue-os/bazzite/*.sixel | shuf -n 1)"'
alias neowofetch='/usr/bin/fastfetch -c /usr/share/ublue-os/bazzite/fastfetch.jsonc --logo "$(ls /usr/share/ublue-os/bazzite/*.sixel | shuf -n 1)"'
alias hyfetch='/usr/bin/fastfetch -c /usr/share/ublue-os/bazzite/fastfetch.jsonc --logo "$(ls /usr/share/ublue-os/bazzite/*.sixel | shuf -n 1)"'
alias fastfetch='/usr/bin/fastfetch -c /usr/share/ublue-os/bazzite/fastfetch.jsonc --logo "$(ls /usr/share/ublue-os/bazzite/*.sixel | shuf -n 1)"'